# systemd/roi.service
[Unit]
Description=ROI Instrument Bridge
After=network.target

[Service]
Type=simple
# Easiest/safest for SocketCAN bring-up since main.py calls sudo ip link set ...
User=root
WorkingDirectory=/opt/roi
EnvironmentFile=-/etc/roi/roi.env
Environment=PYTHONUNBUFFERED=1
# Avoid Rich trying to do fancy terminal control in journald
Environment=TERM=dumb
ExecStart=/opt/roi/.venv/bin/python /opt/roi/main.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
